# client-api/env-scripts/daily.cron
# m h  dom mon dow   command
# used to generate files as JSONs and Shapefiles every days
20 04 * * * /bin/bash /usr/local/exec_daily.sh >> /logs/exec_daily.log 2>&1
