# client-api/env-php/daily.cron
# m h  dom mon dow   command
# Tasks to syncronize and check state of the last syncronize for DETER-B AMZ project
40 06 * * * /usr/local/bin/php /usr/local/php-client/index.php
10 07 * * * /usr/local/bin/php /usr/local/php-client/checkStatus.php
